# OpenWrt/crontab
*/1 * * * * (/root/checkReadUSB1.sh >> /dev/null)
*/1 * * * * (/root/checkReadUSB.sh >> /dev/null)
